# tools/swarm/systemd/hermes-project-onboarding.timer
[Unit]
Description=Advance one pending UAP project onboarding request

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min
RandomizedDelaySec=5s
Persistent=true
Unit=hermes-project-onboarding.service

[Install]
WantedBy=timers.target
